# server/supervisor.conf
[unix_http_server]
file = /home/www/python/shared/run/supervisor.sock
chmod = 0777

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/home/www/python/shared/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

[supervisorctl]
serverurl = unix:///home/www/python/shared/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:www]
command = /home/www/python/current/server/gunicorn.sh
user = www
stdout_logfile = /home/www/python/shared/log/gunicorn_supervisor.log
redirect_stderr = true

[program:celery]
command = /home/www/python/current/server/celery.sh
user = www
stdout_logfile = /home/www/python/shared/log/celery_supervisor.log
redirect_stderr = true
